# sys/.config/systemd/user/org.gnome.GPaste.service.d/override.conf
[Unit]
Requires=x11.service

[Service]
# Try to restart no matter how the process ended
Restart=always
# Restart in 5 second intervals
RestartSec=5
# Try to restart up to 3 times
StartLimitBurst=3
