# systemd/.config/systemd/userexample/chapod.timer
[Timer]
OnCalendar=*-*-* 01:00:00

[Install]
WantedBy=mystuff.target
