# systemd/.config/systemd/userexample/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/run/user/%U/dbus/user_bus_socket
