# systemd/.config/systemd/userexample/lock.target
[Unit]
Description=Lock Screen on first login
Wants=wm.target
AllowIsolate=true

[Install]
Alias=default.target
