# systemd/.config/systemd/userexample/monitor.path
[Path]
PathExists=/dev/fb1

[Install]
WantedBy=mystuff.target
