# systemd/.config/systemd/userexample/monitor.service
[Service]
Environment=DISPLAY=:0
Type=oneshot
ExecStart=/home/daniel/bin/monitors
RemainAfterExit=yes
