# systemd/.config/systemd/userexample/mystuff.target
[Unit]
Description=Xinitrc Stuff
After=xorg.target

[Install]
Alias=default.target
