# systemd/.config/systemd/userexample/thunderbird.service
[Unit]
Description=Thunderbird client

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/thunderbird

[Install]
WantedBy=mystuff.target
