# systemd/.config/systemd/userexample/wm.target
[Unit]
Description=Window Manager
Wants=xorg.target
Wants=mystuff.target
Requires=dbus.socket
AllowIsolate=true

[Install]
Alias=default.target
