# systemd/.config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/bin/dash %h/.config/mutt/sync_mail
