# systemd/.config/systemd/user/offlineimap.timer
[Timer]
OnUnitInactiveSec=120s
Unit=offlineimap.service
