# systemd/.config/systemd/user/sxhkd.service
[Unit]
Description=hotkey manager
After=dbus.service

[Service]
Type=simple
Restart=on-failure
RestartSec=4
Environment="SXHKD_SHELL=zsh"
ExecStart=/usr/bin/sxhkd
ExecReload=/usr/bin/kill -SIGUSR1 $MAINPID

[Install]
WantedBy=default.target

