# systemd/.config/systemd/user/sxhkd.service
[Unit]
Description=hotkey manager
After=dbus.service

[Service]
Type=simple
Restart=on-failure
RestartSec=2
ExecStart=sxhkd

[Install]
WantedBy=default.target

