# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=handle automounting

[Service]
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target
