# systemd/.config/systemd/user/wallpaper.service
[Unit]
Description=setup wallpaper
After=compton.service

[Service]
Type=forking
ExecStart=%h/bin/wl -r

[Install]
WantedBy=default.target
