# sys/.config/systemd/user/dex.service
[Unit]
Description=Startup with dex
Requires=negwm.service
ConditionPathExists=/usr/bin/dex

[Service]
ExecStart=/usr/bin/dex -a
Restart=on-failure
