# sys/.config/systemd/user/dex.service
[Unit]
Description=Startup with dex
Requires=negwm.service

[Service]
ExecStart=/usr/bin/dex -a
Restart=on-failure
