# sys/.config/systemd/user/dpms-off.service
[Unit]
Description=Disable dpms
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/xset s off -dpms
StartLimitBurst=10
