# sys/.config/systemd/user/isync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/bin/sh %E/mutt/scripts/sync_mail
