# sys/.config/systemd/user/isync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Unit=isync.service

[Install]
WantedBy=timers.target
