# sys/.config/systemd/user/kbd-conf.service
[Unit]
Description=my xinput settings
After=dbus.service

[Service]
Type=forking
ExecStart=%E/xinit/kbd-conf
