# sys/.config/systemd/user/picom.service
[Unit]
Description=Compositing manager
After=dbus.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/picom --backend xrender --dbus --experimental-backends

[Install]
WantedBy=default.target
