# sys/.config/systemd/user/polybar-watcher.path
[Unit]
Description=Enable polybar on change

[Path]
PathChanged=%E/polybar/config.ini

[Install]
WantedBy=graphical-session.target
