# sys/.config/systemd/user/polybar.service
[Unit]
Description=Polybar statusbar
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/polybar -q main
Restart=on-failure
