# sys/.config/systemd/user/polybar.service
[Unit]
Description=Polybar statusbar
PartOf=graphical-session.target
StartLimitIntervalSec=60

[Service]
ExecStart=/usr/bin/polybar -q main
ExecStop=%h/.local/bin/polybar-msg cmd quit
Restart=on-failure
RestartSec=3
StartLimitBurst=30
