# sys/.config/systemd/user/polybar.path
[Path]
PathChanged=/tmp/polybar_activator

[Install]
WantedBy=graphical-session.target
