# sys/.config/systemd/user/udiskie.service
[Unit]
Description=handle automounting
After=dbus.service

[Service]
Type=simple
Restart=always
RestartSec=2
ExecStart=/usr/bin/udiskie -2
Restart=always
LogNamespace=negwm

[Install]
WantedBy=default.target
