# sys/.config/systemd/user/udiskie.service
[Unit]
Description=handle automounting
After=dbus.service

[Service]
Type=simple
Restart=always
RestartSec=2
StartLimitBurst=8
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target
