# sys/.config/systemd/user/wallpaper.service
[Unit]
Description=Setup wallpaper on startup

[Service]
Type=forking
ExecStart=%h/bin/wl -r
Restart=on-failure
LogNamespace=negwm
