# sys/.config/systemd/user/wallpaper.service
[Unit]
Description=Setup wallpaper on startup
After=i3-wm.service

[Service]
Type=forking
ExecStart=%h/bin/wl -r
Restart=on-failure
