# sys/.config/systemd/user/x11.socket
[Unit]
Description=X11 Server Socket
Before=x11.target

[Socket]
ListenStream=/tmp/.X11-unix/X0

[Install]
WantedBy=sockets.target
