# sys/.config/systemd/user/xrdb.service
[Unit]
Description=Merge xrdb settings
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/xrdb -merge "$XDG_CONFIG_HOME/X11/xresources"
