# sys/.config/systemd/user/xrdb.service
[Unit]
Description=Merge xrdb settings
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/xrdb -merge %E/X11/xresources
