# sys/.config/systemd/user/xss-lock.service
[Unit]
Description=Session locker
After=xorg.service

[Service]
ExecStart=/usr/bin/xss-lock --transfer-sleep-lock -- i3lockr -p 8
Restart=on-failure
